# keeperfx/tools/nightlybld/keeperfx_build.cron
#******************************************************************************
#   @file keeperfx_build.cron
#      A 'cron' job to checkout and recompile KeeperFX every night.
#******************************************************************************
42 3 * * * ~/nightly/keeperfx_build.sh 2>&1 > ~/public_html/keeper/nightly/last_log.txt
